# supervisord-obs.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ─────────────────────────────────────────────────────────────────────────────
;  Start order: postgres → redis → prometheus → loki → tempo →
;               backend (after 30 s) → grafana (after 30 s) →
;               promtail (after 35 s) → nginx (after 90 s)
; ─────────────────────────────────────────────────────────────────────────────

; ── PostgreSQL 16 ─────────────────────────────────────────────────────────────
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres
    -D /var/lib/postgresql/16/main
    -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.log
startsecs=5
stopwaitsecs=30

; ── Redis 7 ───────────────────────────────────────────────────────────────────
[program:redis]
command=/usr/bin/redis-server --daemonize no --loglevel notice
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log
startsecs=3

; ── Prometheus ────────────────────────────────────────────────────────────────
[program:prometheus]
command=/usr/local/bin/prometheus
    --config.file=/etc/prometheus/prometheus.yml
    --storage.tsdb.path=/var/lib/prometheus
    --storage.tsdb.retention.time=15d
    --web.enable-remote-write-receiver
    --web.enable-lifecycle
autostart=true
autorestart=true
priority=22
stdout_logfile=/var/log/supervisor/prometheus.log
stderr_logfile=/var/log/supervisor/prometheus.log
startsecs=5

; ── Loki ──────────────────────────────────────────────────────────────────────
[program:loki]
command=/usr/bin/loki -config.file=/etc/loki/loki-config.yml
autostart=true
autorestart=true
priority=23
stdout_logfile=/var/log/supervisor/loki.log
stderr_logfile=/var/log/supervisor/loki.log
startsecs=5

; ── Tempo ─────────────────────────────────────────────────────────────────────
[program:tempo]
command=/usr/bin/tempo -config.file=/etc/tempo/tempo-config.yml
autostart=true
autorestart=true
priority=24
stdout_logfile=/var/log/supervisor/tempo.log
stderr_logfile=/var/log/supervisor/tempo.log
startsecs=5

; ── Spring Boot backend ────────────────────────────────────────────────────────
; Wait 30 s so postgres + redis are ready before Spring Boot connects.
[program:backend]
command=/bin/bash -c "sleep 30 && exec java ${JAVA_OPTS:-} -jar /app/backend.jar"
autostart=true
autorestart=true
priority=30
environment=TEMPO_HOST="localhost",TEMPO_ZIPKIN_PORT="9411"
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.log
startsecs=60
stopwaitsecs=30

; ── Grafana ───────────────────────────────────────────────────────────────────
; Wait 30 s so Prometheus/Loki/Tempo are ready before Grafana UI starts.
[program:grafana]
command=/bin/bash -c "sleep 30 && exec /usr/sbin/grafana-server --homepath=/usr/share/grafana --config=/etc/grafana/grafana.ini web"
autostart=true
autorestart=true
priority=35
stdout_logfile=/var/log/supervisor/grafana.log
stderr_logfile=/var/log/supervisor/grafana.log
startsecs=10

; ── Promtail ──────────────────────────────────────────────────────────────────
; Wait 35 s so Loki is ready and the backend has had a chance to create log files.
[program:promtail]
command=/bin/bash -c "sleep 35 && exec /usr/bin/promtail -config.file=/etc/promtail/promtail-config.yml"
autostart=true
autorestart=true
priority=36
stdout_logfile=/var/log/supervisor/promtail.log
stderr_logfile=/var/log/supervisor/promtail.log
startsecs=5

; ── Nginx (Angular + reverse-proxy) ───────────────────────────────────────────
; Wait 90 s so the backend is up before nginx begins accepting traffic.
[program:nginx]
command=/bin/bash -c "sleep 90 && exec nginx -g 'daemon off;'"
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
startsecs=5
